# alloc-timesheet-bullshit-hourly.service
[Service]
ExecStart=alloc-timesheet-bullshit --systemd-timer
