# alloc-timesheet-bullshit-monthly.service
[Service]
ExecSearchPath=/usr/local/bin/usr/bin:/bin:%h/.local/bin
ExecStart=alloc-timesheet-bullshit --systemd-timer --all
