# nag-laptop.timer
# -*-conf-*-
# Every ten minutes, weekdays only, 10AM to 6PM only.
[Timer]
OnCalendar=Mon..Fri 10..18:00/10
[Install]
WantedBy=timers.target
