# nag-laptop.service
# -*-conf-*-
[Service]
Environment=INSIDE_CRON=true
ExecStart=/home/twb/.bin/nag
