# nag-laptop.service
# -*-conf-*-
[Service]
ExecStart=screen nag-laptop
