# nag-laptop.service
[Service]
ExecSearchPath=/usr/local/bin/usr/bin:/bin:%h/.local/bin
ExecStart=screen nag-laptop
