# nag2al.timer
[Timer]
OnCalendar=*-*-* 05:03:00
[Install]
WantedBy=timers.target
